# trusdx-ai.service
[Unit]
Description=truSDX-AI Radio Driver
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
User=trusdx
Group=trusdx
WorkingDirectory=/opt/trusdx-ai
ExecStart=/usr/bin/python3 /opt/trusdx-ai/src/main.py --verbose --syslog
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/opt/trusdx-ai/src
Environment=HOME=/home/trusdx

# Resource limits
LimitNOFILE=1024
LimitNPROC=100

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/trusdx/.cache/trusdx
ReadWritePaths=/tmp/trusdx_cat
ReadWritePaths=/dev/ttyUSB0
ReadWritePaths=/dev/ttyACM0

[Install]
WantedBy=multi-user.target
